# src/kiran-system-daemon.service.in
[Unit]
Description=System Daemon Service

[Service]
Type=dbus
BusName=com.unikylin.Kiran.SystemDaemon
ExecStart=@EXECDIR@/kiran-system-daemon

